# packages/bdworkbench/bdworkbench-3.2.2.tar.gz/bdworkbench-3.2.2/bdworkbench/baseimg/ubuntu/ubuntu16/supervisord.conf
[supervisord]
user = root
nodaemon = true
logfile = /var/log/supervisord.log

[program:sshd]
command=/bin/bash -c "exec /usr/sbin/sshd -D  > /dev/null 2>&1"
