# packages/window-tracker-x11/window_tracker_x11-0.0.12-py2.py3-none-any.whl/window_tracker_x11-0.0.12.data/data/lib/systemd/user/window-tracker-x11.service
[Unit]
Description=Track X11 window usage into log file
After=display-manager.service

[Service]
ExecStart=/usr/bin/window-tracker-x11
Restart=always

[Install]
WantedBy=default.target
